# crontab
5 12 * * * root /root/job.sh >> /var/log/cron.log 2>&1
